# supervisord.conf
[supervisord]
nodaemon=true

[program:fum]
command: ./manage.py runserver --nostatic 0.0.0.0:8080

[program:postgres]
command: service postgresql restart
